# ansible/roles/mac/files/Library/LaunchDaemons/com.voicevox.engine.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.voicevox.engine</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/share/voicevox/run</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>50021</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/usr/local/share/voicevox</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/voicevox.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/voicevox.error.log</string>
</dict>
</plist>
